# tools/ipvsadm/debian/ipvsadm.files
sbin/ipvsadm
sbin/ipvsadm-restore
sbin/ipvsadm-save
usr/man/man8/ipvsadm.8
etc/init.d/ipvsadm
